# docker/production/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

; PHP-FPM Service
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F -R
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=1
user=appuser

; Laravel Queue Worker (Process background jobs)
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=appuser
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker.log
stopwaitsecs=3600
priority=2

; Laravel Scheduler (Run scheduled tasks every minute)
[program:laravel-scheduler]
command=/bin/bash -c "while true; do /usr/local/bin/php /var/www/html/artisan schedule:run --verbose --no-interaction && sleep 60; done"
autostart=true
autorestart=true
user=appuser
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
priority=3

; Optional: Laravel Horizon (uncomment if using Horizon for queue management)
; [program:laravel-horizon]
; command=/usr/local/bin/php /var/www/html/artisan horizon
; autostart=true
; autorestart=true
; user=appuser
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/horizon.log
; stopwaitsecs=3600
